# dovecot/config/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n 

[program:dovecot]
command=/usr/sbin/dovecot -F

