# fetchmail/config/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n

[program:fetchmail]
command=/usr/bin/fetchmail -f /etc/fetchmailrc

